# usr/lib/systemd/system/systemd-remount-fs.service.d/writableagain.conf
[Service]
# Changing the read/write mode of a bind mount needs a second
# fstab entry, but systemd doesn't support set by itself
ExecStart=/bin/sh -e -c 'if findmnt -O ro -M /etc >/dev/null; then mount -o remount,rw /etc; fi'
